# wiguna.service
[Unit]
Description=Gunicorn instance to serve Wiguna Repair Agent
After=network.target

[Service]
# User running the application
# CHANGE THIS: 'root' might be okay, but a specific user is safer if available
User=root
Group=root

# Root directory of the application
# CHANGE THIS if your path is different from: /home/diagnosa.bengkelwiguna.com/public_html
WorkingDirectory=/home/diagnosa.bengkelwiguna.com/public_html

# Environment variables
# Use the global path and include any other necessary paths
Environment="PATH=/usr/local/bin:/usr/bin:/bin"

# Command to start Gunicorn
ExecStart=/usr/local/bin/gunicorn --workers 3 --bind 127.0.0.1:8000 app:app

[Install]
WantedBy=multi-user.target
